# statistics_for_strava/rootfs/etc/services.d/daemon/run
#!/usr/bin/with-contenv sh
set -eu

timestamp() {
  date -u +"%Y-%m-%dT%H:%M:%SZ"
}

log() {
  echo "$(timestamp) [daemon] $*"
}

log "Starting Statistics for Strava daemon"

RUNTIME_DIR="/data/runtime"
PID_FILE="${RUNTIME_DIR}/daemon.pid"
SOURCE_CONFIG_DIR="/data/config/app"
TARGET_CONFIG_DIR="/var/www/config/app"

mkdir -p "$RUNTIME_DIR"
mkdir -p /var/www/config
mkdir -p "$SOURCE_CONFIG_DIR"

# Reconcile persisted config from add-on options on each service start.
sh /usr/local/bin/sfs-reconcile-config.sh

# Ensure app config path is stable even if the base image recreates /var/www/config.
if [ -L "$TARGET_CONFIG_DIR" ]; then
  rm -f "$TARGET_CONFIG_DIR"
elif [ -d "$TARGET_CONFIG_DIR" ]; then
  rm -rf "$TARGET_CONFIG_DIR"
elif [ -f "$TARGET_CONFIG_DIR" ]; then
  rm -f "$TARGET_CONFIG_DIR"
fi
ln -s "$SOURCE_CONFIG_DIR" "$TARGET_CONFIG_DIR"

# Record daemon PID for health checks; exec preserves this PID.
printf '%s\n' "$$" > "$PID_FILE"

exec bin/console app:daemon:run
