# statistics_for_strava/rootfs/etc/services.d/daemon/run
#!/usr/bin/with-contenv sh
set -eu

timestamp() {
  date +"%Y-%m-%dT%H:%M:%S%z"
}

log() {
  echo "$(timestamp) [daemon] $*"
}

log "Starting Dreeve daemon"

RUNTIME_DIR="/data/runtime"
PID_FILE="${RUNTIME_DIR}/daemon.pid"

mkdir -p "$RUNTIME_DIR"

# Config is reconciled and the /var/www/config/app symlink is created by init
# before this service starts; nothing recreates them afterward. Do NOT reconcile
# or re-link here: reconciling would block the daemon PID (/healthz) behind the
# reconcile lock and re-run migrations every few seconds on a crash loop, and
# re-linking races the web service (both ran `ln -s` → `ln: .../app/app: File exists`).

# Record daemon PID for health checks; exec preserves this PID.
printf '%s\n' "$$" > "$PID_FILE"

exec bin/console app:daemon:run
